# dots/.config/systemd/user/battery-monitor.service
[Unit]
Description=Battery Monitor Service
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/scripts/battery-notify.sh
Environment=DISPLAY=:0
